# src/camconfigd.service
[Unit]
Description=CamFlow configuration service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/camconfigd

[Install]
WantedBy=multi-user.target
